# scripts/systemd/chromexup-daily.timer
[Unit]
Description=Run chromexup daily

[Timer]
Unit=chromexup.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
