# scripts/systemd/chromexup-weekly.timer
[Unit]
Description=Run chromexup weekly

[Timer]
Unit=chromexup.service
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
